# jjj.conf
[program:jjj]
command=python3 /home/curator2/jjj/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/jjj/test.err.log
stdout_logfile=/home/curator2/jjj/test.out.log